# scripts/intra-epoch-bot.service
[Unit]
Description=Intra-Epoch Momentum Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/polymarket-autotrader
ExecStart=/usr/bin/python3 /opt/polymarket-autotrader/bot/intra_epoch_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/polymarket-autotrader/intra_epoch_bot.log
StandardError=append:/opt/polymarket-autotrader/intra_epoch_bot.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
